# deploy/huggingface/all-in-one/supervisord.conf
[supervisord]
nodaemon=true
logfile=/home/user/logs/supervisord.log
pidfile=/home/user/logs/supervisord.pid
user=user
childlogdir=/home/user/logs

[program:redis]
command=redis-server --bind 127.0.0.1 --port 6379 --appendonly no --save "" --dir /home/user/data
autorestart=true
priority=10
stdout_logfile=/home/user/logs/redis.out.log
stderr_logfile=/home/user/logs/redis.err.log

[program:webhook]
command=gh-deepagent serve --host 127.0.0.1 --port 8080
autorestart=true
priority=20
stdout_logfile=/home/user/logs/webhook.out.log
stderr_logfile=/home/user/logs/webhook.err.log
environment=
    DEEPAGENT_REDIS_URL="redis://localhost:6379/0",
    DEEPAGENT_JSON_LOGS="1"

[program:worker]
command=gh-deepagent worker --workers 1
autorestart=true
priority=30
stdout_logfile=/home/user/logs/worker.out.log
stderr_logfile=/home/user/logs/worker.err.log
environment=
    DEEPAGENT_REDIS_URL="redis://localhost:6379/0",
    DEEPAGENT_JSON_LOGS="1"

# Streamlit is the public face — exposed on the Space's only port (7860).
# It talks to the webhook via 127.0.0.1:8080 (set via DEEPAGENT_API_URL in
# the Dockerfile). GitHub cannot call this Space's webhook directly; that's
# fine for a demo.
[program:streamlit]
command=streamlit run src/gh_deepagent/dashboard/app.py
    --server.address=0.0.0.0 --server.port=7860
    --server.headless=true
    --browser.gatherUsageStats=false
autorestart=true
priority=50
stdout_logfile=/home/user/logs/streamlit.out.log
stderr_logfile=/home/user/logs/streamlit.err.log
