# deploy/huggingface/all-in-one/supervisord.conf
[supervisord]
nodaemon=true
logfile=/home/user/logs/supervisord.log
pidfile=/home/user/logs/supervisord.pid
user=user
childlogdir=/home/user/logs

[program:redis]
command=redis-server --bind 127.0.0.1 --port 6379 --appendonly no --save "" --dir /home/user/data
autorestart=true
priority=10
stdout_logfile=/home/user/logs/redis.out.log
stderr_logfile=/home/user/logs/redis.err.log

[program:webhook]
command=gh-deepagent serve --host 0.0.0.0 --port 8080
autorestart=true
priority=20
stdout_logfile=/home/user/logs/webhook.out.log
stderr_logfile=/home/user/logs/webhook.err.log
environment=
    DEEPAGENT_REDIS_URL="redis://localhost:6379/0",
    DEEPAGENT_JSON_LOGS="1"

[program:worker]
command=gh-deepagent worker --workers 1
autorestart=true
priority=30
stdout_logfile=/home/user/logs/worker.out.log
stderr_logfile=/home/user/logs/worker.err.log
environment=
    DEEPAGENT_REDIS_URL="redis://localhost:6379/0",
    DEEPAGENT_JSON_LOGS="1"

# Public-facing router on port 7860. Routes:
#   /webhook       → webhook (so GitHub can call <space>.hf.space/webhook)
#   /healthz       → webhook
#   /metrics       → webhook (admin-only via bearer auth, see SAAS_FEATURES.md)
#   anything else  → Streamlit dashboard
[program:router]
command=python /home/user/app/router.py
autorestart=true
priority=40
stdout_logfile=/home/user/logs/router.out.log
stderr_logfile=/home/user/logs/router.err.log

[program:streamlit]
command=streamlit run src/gh_deepagent/dashboard/app.py
    --server.address=127.0.0.1 --server.port=8501
    --server.enableCORS=false --server.enableXsrfProtection=true
    --server.baseUrlPath=/
autorestart=true
priority=50
stdout_logfile=/home/user/logs/streamlit.out.log
stderr_logfile=/home/user/logs/streamlit.err.log
